# irradium-contrib/hostapd/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/hostapd/
-rw-r--r--	root/root	etc/hostapd/hostapd.accept
-rw-r--r--	root/root	etc/hostapd/hostapd.conf
-rw-r--r--	root/root	etc/hostapd/hostapd.deny
-rw-r--r--	root/root	etc/hostapd/hostapd.eap_user
-rw-r--r--	root/root	etc/hostapd/hostapd.radius_clients
-rw-r--r--	root/root	etc/hostapd/hostapd.sim_db
-rw-r--r--	root/root	etc/hostapd/hostapd.wpa_psk
drwxr-xr-x	root/root	etc/rc.d/
-rwxr-xr-x	root/root	etc/rc.d/hostapd
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/bin/
-rwxr-xr-x	root/root	usr/bin/hostapd_cli
-rwxr-xr-x	root/root	usr/bin/nt_password_hash
drwxr-xr-x	root/root	usr/sbin/
-rwxr-xr-x	root/root	usr/sbin/hostapd
drwxr-xr-x	root/root	usr/share/
drwxr-xr-x	root/root	usr/share/man/
drwxr-xr-x	root/root	usr/share/man/man1/
-rw-r--r--	root/root	usr/share/man/man1/hostapd_cli.1.gz
drwxr-xr-x	root/root	usr/share/man/man8/
-rw-r--r--	root/root	usr/share/man/man8/hostapd.8.gz
