# irradium-contrib/rp-pppoe/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/ppp/
drwxr-xr-x	root/root	etc/ppp/plugins/
-rw-r--r--	root/root	etc/ppp/plugins/README
-rwxr-xr-x	root/root	etc/ppp/plugins/rp-pppoe.so
-rw-r--r--	root/root	etc/ppp/pppoe-server-options
drwxr-xr-x	root/root	etc/rc.d/
-rwxr-xr-x	root/root	etc/rc.d/rp-pppoe
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/sbin/
-rwxr-xr-x	root/root	usr/sbin/pppoe
-rwxr-xr-x	root/root	usr/sbin/pppoe-relay
-rwxr-xr-x	root/root	usr/sbin/pppoe-server
-rwxr-xr-x	root/root	usr/sbin/pppoe-sniff
drwxr-xr-x	root/root	usr/share/
drwxr-xr-x	root/root	usr/share/doc/
drwxr-xr-x	root/root	usr/share/doc/rp-pppoe-4.0/
-rw-r--r--	root/root	usr/share/doc/rp-pppoe-4.0/CHANGES
-rw-r--r--	root/root	usr/share/doc/rp-pppoe-4.0/HOW-TO-CONNECT
-rw-r--r--	root/root	usr/share/doc/rp-pppoe-4.0/KERNEL-MODE-PPPOE
-rw-r--r--	root/root	usr/share/doc/rp-pppoe-4.0/LICENSE
-rw-r--r--	root/root	usr/share/doc/rp-pppoe-4.0/README.md
drwxr-xr-x	root/root	usr/share/man/
drwxr-xr-x	root/root	usr/share/man/man8/
-rw-r--r--	root/root	usr/share/man/man8/pppoe-relay.8.gz
-rw-r--r--	root/root	usr/share/man/man8/pppoe-server.8.gz
-rw-r--r--	root/root	usr/share/man/man8/pppoe-sniff.8.gz
-rw-r--r--	root/root	usr/share/man/man8/pppoe.8.gz
